# dot_config/systemd/user/rclone.timer
[Unit]
Description=Rclone timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=rclone.service

[Install]
WantedBy=default.target
